# packages/genocide/genocide-91-py3-none-any.whl/genocide-91.data/data/genocide/genocide.service
[Unit]
Description=GENOCIDE - @KarimKhanQC reconsider OTP-CR-117/19
After=multi-user.target

[Service]
#Type=forking
KillMode=mixed
#PrivateTmp=true
DynamicUser=True
StateDirectory=genocide
LogsDirectory=genocide
CacheDirectory=genocide
ExecStart=/usr/local/bin/genocided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
